# forge-lamp/configs/apache2/supervisor.conf
[program:apache2]
command=/apache2-service.sh
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
